# tests/src/Hodor/Daemon/ExpectedSupervisordConfig.conf
[program:hodor-superqueuer-default]
command=/usr/bin/env php '{{HODOR_BASE_PATH}}/src/Hodor/Daemon/../../../bin/superqueuer.php' -c '{{HODOR_BASE_PATH}}/tests/src/Hodor/Daemon/ExpectedSupervisordConfig.php'
process_name=%(program_name)s_%(process_num)d
numprocs=1
numprocs_start=0
autorestart=true
autostart=true
startsecs=0
startretries=3
user=apache
stopsignal=TERM
stderr_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
stdout_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.debug.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2

[program:hodor-bufferer-default]
command=/usr/bin/env php '{{HODOR_BASE_PATH}}/src/Hodor/Daemon/../../../bin/buffer-worker.php' -c '{{HODOR_BASE_PATH}}/tests/src/Hodor/Daemon/ExpectedSupervisordConfig.php' -q 'default'
process_name=%(program_name)s_%(process_num)d
numprocs=10
numprocs_start=0
autorestart=true
autostart=true
startsecs=0
startretries=3
user=apache
stopsignal=TERM
stderr_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
stdout_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.debug.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2

[program:hodor-bufferer-special]
command=/usr/bin/env php '{{HODOR_BASE_PATH}}/src/Hodor/Daemon/../../../bin/buffer-worker.php' -c '{{HODOR_BASE_PATH}}/tests/src/Hodor/Daemon/ExpectedSupervisordConfig.php' -q 'special'
process_name=%(program_name)s_%(process_num)d
numprocs=10
numprocs_start=0
autorestart=true
autostart=true
startsecs=0
startretries=3
user=apache
stopsignal=TERM
stderr_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
stdout_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.debug.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2

[program:hodor-worker-default]
command=/usr/bin/env php '{{HODOR_BASE_PATH}}/src/Hodor/Daemon/../../../bin/job-worker.php' -c '{{HODOR_BASE_PATH}}/tests/src/Hodor/Daemon/ExpectedSupervisordConfig.php' -q 'default'
process_name=%(program_name)s_%(process_num)d
numprocs=10
numprocs_start=0
autorestart=true
autostart=true
startsecs=0
startretries=3
user=apache
stopsignal=TERM
stderr_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
stdout_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.debug.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2

[program:hodor-worker-intense]
command=/usr/bin/env php '{{HODOR_BASE_PATH}}/src/Hodor/Daemon/../../../bin/job-worker.php' -c '{{HODOR_BASE_PATH}}/tests/src/Hodor/Daemon/ExpectedSupervisordConfig.php' -q 'intense'
process_name=%(program_name)s_%(process_num)d
numprocs=2
numprocs_start=0
autorestart=true
autostart=true
startsecs=0
startretries=3
user=apache
stopsignal=TERM
stderr_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=2
stdout_logfile=/var/log/hodor/%(program_name)s_%(process_num)d.debug.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=2

